# insider-bot.service
[Unit]
Description=Insider Trading Bot - Automated SEC Form 4 Analysis
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User={{BOT_USER}}
Group={{BOT_GROUP}}
WorkingDirectory={{BOT_DIR}}
Environment=PATH={{BOT_DIR}}/insider_bot_env/bin
ExecStart={{BOT_DIR}}/insider_bot_env/bin/python main.py
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutSec=30

# Restart configuration
Restart=always
RestartSec=10
RestartPreventExitStatus=2

# Security settings
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths={{BOT_DIR}}
PrivateTmp=yes

# Logging
StandardOutput=append:{{BOT_DIR}}/logs/service.log
StandardError=append:{{BOT_DIR}}/logs/service.log
SyslogIdentifier=insider-bot

# Resource limits
MemoryLimit=1G
CPUQuota=50%

[Install]
WantedBy=multi-user.target